# deploy/macos/com.marko.clippa.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marko.clippa</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/marko/reclip/scripts/reclip-launch.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/marko/reclip</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOST</key>
        <string>0.0.0.0</string>
        <key>PORT</key>
        <string>8899</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/marko/reclip/logs/clippa.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/marko/reclip/logs/clippa.err.log</string>
</dict>
</plist>
